# 00-系统/运行时/crontab/orbit-worklog.cron
# OrbitOS daily worklog bootstrap.
# Source of truth: 00-系统/运行时/crontab/OrbitOS-worklog.cron
# Installed by: orbit-vault install-runtime --crontab
0 8 * * * /bin/zsh -lc 'node "/Users/zyongzhu/OrbitOS/00-系统/Skills/orbit-vault/scripts/orbit-vault.mjs" ensure-worklog --vault "/Users/zyongzhu/OrbitOS" >/tmp/OrbitOS-worklog.log 2>&1'
